# configs/backup.service
[Unit]
Description=Daily server backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/ubuntu/scripts/backup.sh
User=root
StandardOutput=journal
StandardError=journal
